# dpkg/info/arpon.list
/.
/etc
/etc/arpon.conf
/etc/default
/etc/default/arpon
/etc/init.d
/etc/init.d/arpon
/etc/logrotate.d
/etc/logrotate.d/arpon
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/arpon.service
/lib/systemd/system/arpon@.service
/usr
/usr/sbin
/usr/sbin/arpon
/usr/share
/usr/share/doc
/usr/share/doc/arpon
/usr/share/doc/arpon/AUTHOR
/usr/share/doc/arpon/NEWS.Debian.gz
/usr/share/doc/arpon/README.Debian
/usr/share/doc/arpon/changelog.Debian.gz
/usr/share/doc/arpon/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/arpon.8.gz
