# dpkg/info/cacti.list
/etc
/etc/cacti
/etc/cron.d
/etc/cron.d/cacti
/etc/logrotate.d
/etc/logrotate.d/cacti
/var
/var/log
/var/log/cacti
