# dpkg/info/epson-printer-utility.list
/.
/opt
/opt/epson-printer-utility
/opt/epson-printer-utility/bin
/opt/epson-printer-utility/bin/epson-printer-utility
/opt/epson-printer-utility/lib
/opt/epson-printer-utility/lib/libEPSCommonLib.a
/opt/epson-printer-utility/resource
/opt/epson-printer-utility/resource/Languages
/opt/epson-printer-utility/resource/Languages/utilitytr_ja.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_fr.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_ko.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_pt.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_ru.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_zh_TW.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_zh_CN.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_es.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_nl.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_it.qm
/opt/epson-printer-utility/resource/Languages/utilitytr_de.qm
/opt/epson-printer-utility/resource/Images
/opt/epson-printer-utility/resource/Images/STM.png
/opt/epson-printer-utility/resource/Images/printer_status.png
/opt/epson-printer-utility/resource/Images/Nozzle.png
/opt/epson-printer-utility/resource/Images/warning.png
/opt/epson-printer-utility/resource/Images/status_indicator_offline.png
/opt/epson-printer-utility/resource/Images/Cleaning.png
/opt/epson-printer-utility/resource/Images/printer_status_error_laser.png
/opt/epson-printer-utility/resource/Images/Cleaning_small.png
/opt/epson-printer-utility/resource/Images/printer_status_OLD.png
/opt/epson-printer-utility/resource/Images/Inkout.png
/opt/epson-printer-utility/resource/Images/printer_status_laser.png
/opt/epson-printer-utility/resource/Images/Gap.png
/opt/epson-printer-utility/resource/Images/printer_status_error.png
/opt/epson-printer-utility/resource/Images/LFP2.png
/opt/epson-printer-utility/resource/Images/Cleaning_small_old.png
/opt/epson-printer-utility/resource/Images/Nozzle_small_old.png
/opt/epson-printer-utility/resource/Images/Nozzle_small.png
/opt/epson-printer-utility/resource/Images/printer_status_warning.png
/opt/epson-printer-utility/resource/Images/Inklow.png
/opt/epson-printer-utility/resource/Images/epson.png
/opt/epson-printer-utility/resource/Images/status_indicator_online.png
/opt/epson-printer-utility/resource/Images/Gap_small.png
/opt/epson-printer-utility/resource/Images/printer_status_warning_laser.png
/opt/epson-printer-utility/resource/Images/AppIcon.png
/opt/epson-printer-utility/resource/Images/STM_small.png
/opt/epson-printer-utility/resource/Images/Inklow_.png
/opt/epson-printer-utility/doc
/opt/epson-printer-utility/doc/AUTHORS
/opt/epson-printer-utility/doc/README
/opt/epson-printer-utility/doc/INSTALL
/opt/epson-printer-utility/doc/ChangeLog
/opt/epson-printer-utility/doc/NEWS
/opt/epson-printer-utility/rules
/opt/epson-printer-utility/rules/79-udev-epson.rules
/usr
/usr/lib
/usr/lib/epson-backend
/usr/lib/epson-backend/ecbd
/usr/lib/epson-backend/rc.d
/usr/lib/epson-backend/rc.d/ecbd
/usr/lib/epson-backend/rc.d/init-functions
/usr/lib/epson-backend/rc.d/inst-rc_d.sh
/usr/lib/epson-backend/scripts
/usr/lib/epson-backend/scripts/inst-cups-post.sh
/usr/lib/cups
/usr/lib/cups/backend
/usr/lib/cups/backend/ecblp
