# dpkg/info/hddtemp.conffiles
/etc/hddtemp.db
/etc/init.d/hddtemp
/etc/logcheck/ignore.d.server/hddtemp
